# conf.d/console.conf
[program:console]
command=/root/.startjupyter.sh 9000
directory=/var/www/data/admin/files/Tensorflow-Tutorials
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
auto_start=true
autorestart=true
